# cv研究生日常Lab/更新服务器用户占用情况/getuser.conf
[program:getuser]
command = /home/ywz/anaconda3/bin/python /home/ywz/data/get_user/get_user.py
autostart = true
