# cv研究生日常Lab/frp/supervisor自启动/frp.conf
[program:frp]
command = /root/frp_0.29.0_linux_amd64/frps -c /root/frp_0.29.0_linux_amd64/frps_full.ini
autostart = true