# web相关 (python全栈系列)/vps-frp-proxy/supervisor/sele.conf
[program:sele]
command = python3 /root/down.py
autostart = true

