# web相关 (python全栈系列)/安卓termux变服务器/conf.d/frpc2.conf
[program:frp2]
command = /home/ywz/data/frp/frpc -c /home/ywz/data/frp/frpc2.ini
autostart = true
user = root
