# web相关 (python全栈系列)/安卓termux变服务器/conf.d/frpc.conf
[program:frp]
command = /home/ywz/data/frp/frpc -c /home/ywz/data/frp/frpc.ini
autostart = true
user = root
