# web相关 (python全栈系列)/nginx+uwsgi部署django/F46059.conf
[program:F46059]
command=/home/sharklet/anaconda3/bin/python -u manage.py runserver 0.0.0.0:8001
directory=/var/www/html/F46059/django1/
stopasgroup=true
stopsignal=QUIT
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/F46059/django1/out.log
stderr_logfile=/var/www/html/F46059/django1/err.log
loglevel=info
