# web相关 (python全栈系列)/supervisor管理django/01181.conf
[program:01181]
command=/usr/bin/python3 -u manage.py runserver 0.0.0.0:8001
directory=/var/www/clients/01181/django1/
stopasgroup=true
stopsignal=QUIT
user=pi
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/clients/01181/django1/01181_out.log
stderr_logfile=/var/www/clients/01181/django1/01181_err.log
loglevel=info
