# cv研究生日常Lab/frp/frp2.conf
[program:frp]
command = frpc -c frpc.ini
directory = /home/sharklet/data/frp_0.29.0_linux_amd64/
autostart = true
user = root