# cv研究生日常Lab/pySrun4k_BeihangLogin-master/pylogin.conf
[program:pylogin]
command = /home/ywz/anaconda3/bin/python /home/ywz/data/pySrun4k_BeihangLogin-master/main_login_regular.py
autostart = true
